# netbird/rootfs/etc/s6-overlay/s6-rc.d/netbird/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: NetBird Client
# Runs NetBird Client
# ==============================================================================
declare -a options
declare name
declare value


# Get the options configured in HASS GUI
readonly CONFIG_OLD_PATH=/homeassistant/netbird/config.json
readonly CONFIG_PATH=/var/lib/netbird/config.json
readonly CONFIG_PROFILE_PATH=/var/lib/netbird/default.json

# Migrate old configuration files to new location (one-time migration)
if [ -f "${CONFIG_OLD_PATH}" ] && [ ! -f "${CONFIG_PROFILE_PATH}" ]; then
    bashio::log.info "Migrating configuration from old location..."
    mv "${CONFIG_OLD_PATH}" "${CONFIG_PATH}" && bashio::log.info "Migration step 1 complete" || bashio::log.error "Migration step 1 failed: could not move ${CONFIG_OLD_PATH} to ${CONFIG_PATH}"
fi

if [ -f "${CONFIG_PATH}" ] && [ ! -f "${CONFIG_PROFILE_PATH}" ]; then
    bashio::log.info "Setting up NetBird profile configuration..."
    mv "${CONFIG_PATH}" "${CONFIG_PROFILE_PATH}" && bashio::log.info "Migration step 2 complete"
fi

admin_url="$(bashio::config 'admin_url')"
management_url="$(bashio::config 'management_url')"
setup_key="$(bashio::config 'setup_key')"
hostname="$(bashio::config 'hostname')"
rosenpass="$(bashio::config 'rosenpass')"
rosenpass_permissive="$(bashio::config 'rosenpass_permissive')"
log_level="$(bashio::config 'log_level')"

options+=(--foreground-mode)
options+=(--log-file console)

if [ "${admin_url}" = "" ]; then
    bashio::log.info "Using Default Admin URL"
else
    bashio::log.info "Using ${admin_url} as Admin URL"
    options+=(--admin-url "${admin_url}")
fi

if [ "${management_url}" = "" ]; then
    bashio::log.info "Using Default Management URL"
else
    bashio::log.info "Using ${management_url} as Management URL"
    options+=(--management-url "${management_url}")
fi

if [ "${setup_key}" = "" ]; then
    bashio::log.info "No Setup Key Set"
    bashio::log.info "This client will only show up in dashboards it's already registered with."
else
    bashio::log.info "Setup Key configured (hidden for security)"
    options+=(--setup-key "${setup_key}")
fi

if [ "${hostname}" = "" ]; then
    bashio::log.info "No Hostname Set"
    bashio::log.info "This client will use the default (<docker container id>-netbird-client) as hostname in peers."
else
    bashio::log.info "Using ${hostname} as hostname"
    options+=(--hostname "${hostname}")
fi

if ! ${rosenpass}; then
    bashio::log.info "Rosenpass disabled"
    options+=(--enable-rosenpass=false)
else
    bashio::log.info "Rosenpass enabled"
    options+=(--enable-rosenpass)
    if ${rosenpass_permissive}; then
        bashio::log.info "Rosenpass permissive mode enabled"
        options+=(--rosenpass-permissive)
    fi
fi

if [ "${log_level}" = "" ] || [ "${log_level}" = "null" ]; then
    bashio::log.info "No log level Set"
    bashio::log.info "This client will use the default logging."
else
    bashio::log.info "Using ${log_level} as logging level"
    options+=(--log-level "${log_level}")
fi

# Load custom environment variables
for var in $(bashio::config 'env_vars|keys'); do
    name=$(bashio::config "env_vars[${var}].name")
    value=$(bashio::config "env_vars[${var}].value")
    bashio::log.info "Setting ${name} to ${value}"
    export "${name}=${value}"
done

# Workaround for DNS resolution with systemd-resolved
# NetBird checks for systemd-resolved by looking for a specific comment in /etc/resolv.conf
# This ensures NetBird can properly detect and configure DNS settings on the host
# See: https://github.com/netbirdio/netbird/issues/dns-detection
CONTENT=$(cat /etc/resolv.conf)
echo '# systemd-resolved' > /etc/resolv.conf
echo "$CONTENT" >> /etc/resolv.conf

# Verify netbird binary exists
if ! command -v netbird &> /dev/null; then
    bashio::log.fatal "NetBird binary not found!"
    exit 1
fi

bashio::log.info "Starting NetBird Client..."
bashio::log.info "netbird up " "${options[@]}"
netbird up "${options[@]}"
